# extensions/service-hostapd/extension/usr/local/etc/init.d/services/service-hostapd/finish
#!/bin/sh

exec 2>&1

grep -q "ap0:" /proc/net/dev && iw dev ap0 del
