# extensions/services/serial/run
#!/bin/sh

set -eu

exec 2>&1

echo "Starting"

exec /sbin/getty -L /dev/ttyAMA0 115200 vt100
